# server/packages/callpy/manager.ini
[program:http]
command= python -m http.server
numprocs = 1